# crontab
* * * * * DISPLAY=:0.0 feh --bg-scale --randomize /home/brad/Dropbox/Wallpapers/Arch
