# crontab
* * * * * DISPLAY=:0.0 feh --bg-scale --randomize ~/Dropbox/Wallpapers/Pack1
