# crontab
* * * * * DISPLAY=:0.0 ~/bin/switch_wallpaper.sh
